# backend/cnc_backend.service
# cnc_backend.service  --  Servicio systemd para el Backend FastAPI (UV)
# =======================================================================
# Instalacion en la VM publica (vm-iot-front):
#   sudo cp backend/cnc_backend.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable cnc_backend
#   sudo systemctl start cnc_backend
#
# Ver logs:
#   sudo journalctl -u cnc_backend -f

[Unit]
Description=FLUX CNC -- Backend FastAPI (GET /datos/ / GET /datos/descargar/)
After=network.target mqtt_bridge.service
Wants=mqtt_bridge.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/cnc-iot-backend/backend
EnvironmentFile=/home/ubuntu/cnc-iot-backend/backend/.env
ExecStart=/usr/local/bin/uv run uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
